# deployments/podman-quadlets/dev/auth-service.container
# =========================================================================
# Podman Quadlet - Auth Service (Development)
# =========================================================================
# Systemd unit file for Auth Service in development environment
# Place in ~/.config/containers/systemd/ or /etc/containers/systemd/

[Unit]
Description=Reciprocal Clubs Auth Service (Development)
Documentation=https://docs.reciprocal-clubs.com/auth-service
Wants=network-online.target postgres-dev.service nats-dev.service redis-dev.service hanko-dev.service
After=network-online.target postgres-dev.service nats-dev.service redis-dev.service hanko-dev.service
RequiresMountsFor=%t/containers

[Container]
Image=reciprocal-clubs/auth-service:dev
ContainerName=reciprocal-auth-service-dev
AutoUpdate=registry

# Development Environment Variables
Environment=APP_ENV=development
Environment=AUTH_SERVICE_SERVICE_PORT=8081
Environment=AUTH_SERVICE_SERVICE_GRPC_PORT=9081
Environment=AUTH_SERVICE_DATABASE_HOST=reciprocal-postgres-dev
Environment=AUTH_SERVICE_DATABASE_PORT=5432
Environment=AUTH_SERVICE_DATABASE_USERNAME=postgres
Environment=AUTH_SERVICE_DATABASE_PASSWORD=devpassword
Environment=AUTH_SERVICE_DATABASE_DATABASE=auth_service
Environment=AUTH_SERVICE_NATS_URL=nats://reciprocal-nats-dev:4222
Environment=AUTH_SERVICE_REDIS_HOST=reciprocal-redis-dev
Environment=AUTH_SERVICE_REDIS_PORT=6379
Environment=AUTH_SERVICE_AUTH_JWT_SECRET=development-secret-key-change-in-production
Environment=AUTH_SERVICE_HANKO_BASE_URL=http://reciprocal-hanko-dev:8000
Environment=AUTH_SERVICE_HANKO_API_KEY=
Environment=AUTH_SERVICE_MFA_ISSUER=Reciprocal Clubs Dev
Environment=AUTH_SERVICE_DEBUG=true
Environment=AUTH_SERVICE_LOG_LEVEL=debug

# Network and Ports
Network=reciprocal-clubs-dev.network
PublishPort=8081:8081
PublishPort=9081:9081

# Health Check
HealthCmd=curl -f http://localhost:8081/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s

# Resource Limits (Development)
Memory=1g
MemorySwap=2g
CPUs=1.0

# Volume Mounts (Development)
Volume=/var/log/reciprocal-clubs/auth-service:/app/logs:Z
Volume=/etc/reciprocal-clubs/dev:/app/config:ro,Z

# Security (Relaxed for Development)
SecurityOpt=label=type:container_runtime_t
User=1001:1001

# Development specific
Label=environment=development
Label=service=auth-service
Label=version=dev
Label=maintainer=reciprocal-clubs-team

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=900
TimeoutStopSec=70

# Development specific service configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auth-service-dev

[Install]
WantedBy=multi-user.target default.target